# nastkom-dcp2/etc/systemd/system/media-sd.automount
[Unit]
Description=sd card mount

[Automount]
Where=/media/sd

[Install]
WantedBy=multi-user.target
